# twib/twibd/twibd.service
[Unit]
Description=Twili Bridge Daemon
Requires=twibd.socket

[Service]
Type=notify
NotifyAccess=all
ExecStart=twibd -s
StandardError=journal
